# tools/data/systemd/mcp-data.service
[Unit]
# Hardened systemd SYSTEM unit for an mcp-tools server (mcp-data). SYSTEM unit (not
# --user) so (a) IPAddressDeny/Allow actually enforce the egress wall and (b) it
# starts at boot. Install via scripts/install-system.sh (sudo once).
Description=MCP tool: mcp-data (public, Google-OAuth gated via Cloudflare Tunnel)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=wes
Group=wes
WorkingDirectory=/home/wes/Documents/mcp-tools/tools/data
# MCP_PORT is read by main() BEFORE .env loads, so it MUST live here, not in .env.
Environment=MCP_HOST=127.0.0.1
Environment=MCP_PORT=8062
Environment=PATH=/home/wes/.local/bin:/usr/local/bin:/usr/bin:/bin
# OAuth/FastMCP state AND OpenBB's $HOME/.openbb_platform settings stay writable under
# ProtectHome=read-only via the StateDirectory (/var/lib/mcp-data). HOME points there so
# OpenBB (which derives its settings dir from Path.home()) writes inside it, not the
# read-only real home. The data lake lives there too (DATA_ROOT below).
StateDirectory=mcp-data
Environment=HOME=%S/mcp-data
Environment=FASTMCP_HOME=%S/mcp-data
Environment=DATA_ROOT=%S/mcp-data/data
# Freeze OpenBB's prebuilt accessor: never rebuild (i.e. never write into the read-only
# venv) at import time. Prebuild at setup with `.venv/bin/python -c "import openbb"`.
Environment=OPENBB_AUTO_BUILD=false

# --- L2 egress: all outbound traffic MUST go through the egress proxy ---
# This tool's per-tool squid listener (:8074) enforces its domain
# allowlist (security/egress-proxy/allowlist/data.txt). The IPAddress wall makes
# the proxy the ONLY route off-box; NO_PROXY keeps loopback (guardrail, resolver) direct.
Environment=HTTPS_PROXY=http://127.0.0.1:8074
Environment=HTTP_PROXY=http://127.0.0.1:8074
Environment=NO_PROXY=127.0.0.1,localhost
IPAddressDeny=any
IPAddressAllow=127.0.0.1/8

ExecStart=/home/wes/Documents/mcp-tools/tools/data/.venv/bin/python server.py
Restart=on-failure
RestartSec=5

# --- Hardening (L3) ---
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=read-only
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX
RestrictNamespaces=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
LockPersonality=yes
SystemCallArchitectures=native
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM
# Enforceable as a SYSTEM unit (these silently failed in --user mode):
CapabilityBoundingSet=
AmbientCapabilities=
PrivateDevices=yes
ProtectClock=yes
ProtectHostname=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectKernelLogs=yes
ProtectControlGroups=yes
ProtectProc=invisible
ProcSubset=pid

[Install]
WantedBy=multi-user.target
